# deploy/systemd/fitness-food-refresh.timer
[Unit]
Description=Hourly pending-check refresh for Fitness Food materialized views

[Timer]
OnCalendar=hourly
Persistent=true
Unit=fitness-food-refresh.service

[Install]
WantedBy=timers.target
