# deploy/systemd/fitness-food-refresh.timer
[Unit]
Description=Daily refresh for Fitness Food materialized views

[Timer]
OnCalendar=*-*-* 04:10:00
Persistent=true
Unit=fitness-food-refresh.service

[Install]
WantedBy=timers.target
